# samba/rootfs/etc/s6-overlay/s6-rc.d/mdnsd/finish
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash

killall mdnsd
